# crontab
50 1 * * *  cd ~/Desktop/Kartees_Data_Collection; git add .; git commit -m "next set of conso"; git push
